# systemd/system/fan-control.service
[Unit]
Description=AIO Fan Control
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fan-control.sh
